# buttonosbase/etc/udev/rules.d/76-wicd.rules
ACTION!="add", GOTO "wicd_end"
SUBSYSTEM!="net", GOTO "wicd_end"

RUN+="/usr/lib/buttonos/wicdnetif %k"
LABEL="wicd_end"
